# ws/jammy/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:workstation]
command = /google/scripts/entrypoint.sh
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:dockerd]
command=/usr/bin/dockerd
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0